# claude-life-agent.service
[Unit]
Description=Claude Life Agent
Documentation=file:/opt/claude-life-agent/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/opt/claude-life-agent
EnvironmentFile=/opt/claude-life-agent/.env
ExecStart=/opt/claude-life-agent/venv/bin/python main.py
Restart=on-failure
RestartSec=10
SyslogIdentifier=claude-life-agent

[Install]
WantedBy=multi-user.target
